# x11vnc-auto.socket
[Unit]
Description=x11vnc Socket for Per-Connection Servers

[Socket]
ListenStream=0.0.0.0:5900
Accept=yes

[Install]
WantedBy=sockets.target
